# ytdl-sync.service
[Unit]
Description=Youtube-dl downloader service

[Service]
Type=oneshot
ExecStart=ytdl-sync
TimeoutStartSec=3h
