# ytdl-sync.service
[Unit]
Description=Youtube-dl downloader service

[Service]
Type=simple
ExecStart=ytdl-sync
RuntimeMaxSec=3h
