# ytdl-sync.timer
[Unit]
Description=Youtube-dl timer
After=network.target

[Timer]
# Daily
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=multi-user.target
